# sautExternal/board/naqada_panel/overlay/etc/udev/rules.d/10-wi-fi-stop.rules
ACTION=="remove", SUBSYSTEM=="usb", ENV{ID_VENDOR_ID}=="0bda", ENV{ID_MODEL_ID}=="f179", RUN+="/bin/sh -c 'killall hostapd > /dev/null'"
